# MELA-AnalysisService/MELA-ElasticityDependencyAnalysisService-tar/scripts/mela-dependency-analysis-service
### BEGIN INIT INFO
# Provides:          mela-dependency-analysis-service
# Required-Start: mysqld
# Required-Stop:
# Default-Start:     
# Default-Stop:      0 1 6
# Short-Description: Start MELA ElasticityRelationshipsAnalysisService daemon
# Description:       To start-stop MELA-DS
### END INIT INFO

RUN_MODE="daemons"
. /etc/profile
JAVA=$(which java)
NAME=mela-dependency-analysis-service
DAEMONDIR=/opt/mela-dependency-analysis-service
DAEMON=MELA-ElasticityRelationshipsAnalysisService-1.0-SNAPSHOT-exec-war.jar
PIDDIR=/tmp/
PIDFILE=$PIDDIR/mela-dependency-analysis-service.pid
LOG_DIR=/tmp/mela
MELA_HTTP_PORT=8183
PROFILES=MYSQL
#PROFILES=H2
#PROFILES=PORTABLE

test -x $JAVA -a -f $DAEMONDIR/$DAEMON || echo "test $DAEMONDIR/$DAEMON dod not succed" 
test -x $JAVA -a -f $DAEMONDIR/$DAEMON || exit 0
 

. /etc/init.d/functions

case "$1" in
	start)
		echo "Starting Mela-ElasticityRelationshipsAnalysisService daemon"
		mkdir -p $LOG_DIR
                cd $DAEMONDIR
                daemon --pidfile=$PIDFILE $JAVA -DLOG_DIR=$LOG_DIR -jar $DAEMON -httpPort $MELA_HTTP_PORT > $LOG_DIR/mela-dependency-analysis.log -Dspring.profiles.active=$PROFILES 2>&1 &
		RETVAL=$?
                echo
		[ $RETVAL = 0 ] && touch /var/lock/mela-dependency-analysis-service-daemon
		;;
	stop)
		echo "Stopping Mela-ElasticityRelationshipsAnalysisService daemon"
		#killproc --pidfile $PIDFILE
		pid=`ps -ef | grep $DAEMON | tr -s ' '  | cut -f2 -d ' ' | head -n 2`
                echo $pid
                sudo kill -9 $pid
		RETVAL=$?
   		echo
                [ $RETVAL = 0 ] && rm -f /var/lock/mela-dependency-analysis-service-daemon
 		;;
	restart)
		$0 stop
		sleep 1
		$0 start
		;;
	status)
		pidofproc -p $PIDFILE $JAVA >/dev/null
		status=$?
		if [ $status -eq 0 ]; then
			echo "Mela-ElasticityRelationshipsAnalysisService is running"
		else
			echo "Mela-ElasticityRelationshipsAnalysisService is not running"
		fi
		exit $status
		;;
	*)
		echo "Usage: $0 {start|stop|restart|status}"
		exit 1
		;;
esac

exit 0

